# mkosi.extra/etc/systemd/system/nvidia-fabricmanager.service.d/override.conf
[Unit]

[Service]
ExecCondition=/usr/local/bin/check-nvidia-fabric
