# mkosi.extra/etc/systemd/system/tinfoil-egress.service
[Unit]
Description=Tinfoil Container Egress Filter Service
After=nftables.service
Requires=nftables.service tinfoil-boot.service

[Service]
# Type=notify so `systemctl start` blocks until the daemon signals READY=1
# after the first successful DNS resolution and nftables update. Any failure
# in the initial population surfaces as a boot error in tinfoil-boot.
Type=notify
ExecStart=/usr/bin/tinfoil-egress
Restart=on-failure
RestartSec=5s
User=root
Group=root
StandardOutput=journal+console
StandardError=journal+console
