# mkosi.extra/etc/systemd/system/nvidia-persistenced.service.d/override.conf
[Unit]
After=nvidia-fabricmanager.service

[Service]
ExecCondition=/usr/local/bin/check-nvidia-gpu
ExecStartPre=/sbin/modprobe nvidia
ExecStartPre=/sbin/modprobe nvidia_uvm
ExecStart=
ExecStart=/usr/bin/nvidia-persistenced --user nvidia-persistenced --uvm-persistence-mode --verbose
TimeoutSec=1800
